# docker/web/supervisor.ini
[program:uwsgi]
command = uwsgi --ini /app/uwsgi.ini

[program:nginx]
command = /usr/sbin/nginx